# root/etc/services.d/entr/run
#!/usr/bin/with-contenv sh

s6-setuidgid abc while true; do find $FOLDER | entr -d rclone sync; done
