# root/etc/services.d/rcd/run
#!/usr/bin/with-contenv sh

s6-setuidgid abc rclone rcd -q --rc-addr=0.0.0.0:5572 ${RC_USER:+--rc-user=$RC_USER} ${RC_PASS:+--rc-pass=$RC_PASS}
